# samson/debian/prerm
#!/bin/bash
/etc/init.d/samson stop
update-rc.d -f samson remove > /dev/null
userdel samson
