# cosmos/platform/frontend/cosmos-frontend.cron
#!/bin/bash
export PATH=/usr/local/cosmos-frontend/bin:/bin:/usr/bin

django update_status --quiet
sleep 20
django update_status --quiet
sleep 20
django update_status --quiet
