# etc/profile.d/samson.sh
# Environment script for the SAMSON platform
SAMSON_HOME=/opt/samson
SAMSON_WORKING=/var/samson
PATH=${PATH}:${SAMSON_HOME}/bin

export SAMSON_HOME SAMSON_WORKING PATH
